# templates/ansible-docker-update.timer.j2
[Unit]
Description=Update docker-compose images daily

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target